# deploy/macos/com.detec.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.detec.agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>detec-agent</string>
		<string>--interval</string>
		<string>300</string>
		<string>--api-url</string>
		<string>https://api.example.com</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>Library/Logs/DetecAgent/detec-agent.log</string>
	<key>StandardErrorPath</key>
	<string>Library/Logs/DetecAgent/detec-agent.err</string>
	<!-- Do not put API key here. Set AGENTIC_GOV_API_KEY via launchctl setenv or use keychain (see DEPLOY.md). -->
	<!-- Optional: add key EnvironmentVariables with dict AGENTIC_GOV_API_URL, AGENTIC_GOV_INTERVAL if you prefer env over args. -->
</dict>
</plist>
