# packaging/macos/scripts/postinstall
#!/bin/bash
# Post-install script for Detec Agent .pkg
#
# Sets up the LaunchAgent so the agent starts automatically at login.
# Creates the required directory structure and log directory.

set -e

PLIST_LABEL="com.detec.agent"
APP_PATH="/Applications/Detec Agent.app"
BINARY_PATH="${APP_PATH}/Contents/MacOS/detec-agent-gui"
LAUNCH_AGENTS_DIR="$HOME/Library/LaunchAgents"
PLIST_PATH="${LAUNCH_AGENTS_DIR}/${PLIST_LABEL}.plist"
LOG_DIR="$HOME/Library/Logs/DetecAgent"
STATE_DIR="$HOME/.agentic-gov"
APP_SUPPORT_DIR="$HOME/Library/Application Support/Detec"

# Resolve the real (non-root) user who owns $HOME.
# Pkg install scripts run as root, so we need to chown everything back.
REAL_USER="$(stat -f '%Su' "$HOME")"
REAL_GROUP="$(stat -f '%Sg' "$HOME")"

echo "Detec Agent: post-install starting (user=$REAL_USER)..."

# Create required directories
mkdir -p "$LAUNCH_AGENTS_DIR"
mkdir -p "$LOG_DIR"
mkdir -p "$STATE_DIR"
mkdir -p "$APP_SUPPORT_DIR"

# Write the LaunchAgent plist
cat > "$PLIST_PATH" << PLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>${PLIST_LABEL}</string>

    <key>ProgramArguments</key>
    <array>
        <string>${BINARY_PATH}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>${LOG_DIR}/agent.log</string>

    <key>StandardErrorPath</key>
    <string>${LOG_DIR}/agent-error.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
PLIST

echo "  LaunchAgent plist written to $PLIST_PATH"

# Set correct permissions and transfer ownership from root to the real user.
# launchctl requires user-domain plists to be owned by the user.
chmod 644 "$PLIST_PATH"
chmod 755 "$STATE_DIR"
chown "$REAL_USER:$REAL_GROUP" "$PLIST_PATH"
chown -R "$REAL_USER:$REAL_GROUP" "$LOG_DIR"
chown -R "$REAL_USER:$REAL_GROUP" "$STATE_DIR"
chown -R "$REAL_USER:$REAL_GROUP" "$APP_SUPPORT_DIR"

# Copy baked-in config from the .app bundle if present.
BAKED_ENV="${APP_PATH}/Contents/Resources/config/agent.env"
if [ -f "$BAKED_ENV" ]; then
    cp "$BAKED_ENV" "$APP_SUPPORT_DIR/agent.env"
    chmod 600 "$APP_SUPPORT_DIR/agent.env"
    chown "$REAL_USER:$REAL_GROUP" "$APP_SUPPORT_DIR/agent.env"
    echo "  Pre-configured agent.env installed to $APP_SUPPORT_DIR/"
fi

# Load the LaunchAgent as the real user (not root).
REAL_UID="$(id -u "$REAL_USER")"
echo "  Loading LaunchAgent for UID $REAL_UID..."
launchctl bootstrap "gui/$REAL_UID" "$PLIST_PATH" 2>/dev/null || \
    sudo -u "$REAL_USER" launchctl load "$PLIST_PATH" 2>/dev/null || true

echo "Detec Agent: post-install complete."
if [ -f "$BAKED_ENV" ]; then
    echo "  Server config was pre-loaded. The agent will connect automatically."
else
    echo "  The agent will start automatically at login."
    echo "  Run 'detec-agent setup' to configure the API URL and key."
fi
exit 0
