# packaging/macos/scripts/postinstall
#!/bin/bash
# Post-install script for Detec Agent .pkg
#
# Installs a LaunchDaemon so the agent starts at boot (system-wide,
# runs as root). Config is read from /Library/Application Support/Detec.
# Creates /Library/Application Support/Detec and /Library/Logs/DetecAgent.

set -e

PLIST_LABEL="com.detec.agent"
APP_PATH="/Applications/Detec Agent.app"
BINARY_PATH="${APP_PATH}/Contents/MacOS/detec-agent-gui"
LAUNCH_DAEMONS_DIR="/Library/LaunchDaemons"
PLIST_PATH="${LAUNCH_DAEMONS_DIR}/${PLIST_LABEL}.plist"
LOG_DIR="/Library/Logs/DetecAgent"
APP_SUPPORT_DIR="/Library/Application Support/Detec"

echo "Detec Agent: post-install starting (LaunchDaemon, system-wide)..."

# Create system directories (root-owned)
mkdir -p "$APP_SUPPORT_DIR"
mkdir -p "$LOG_DIR"

# Write LaunchDaemon plist (content matches packaging/macos/com.detec.agent.plist)
cat > "$PLIST_PATH" << PLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>${PLIST_LABEL}</string>
    <key>ProgramArguments</key>
    <array>
        <string>${BINARY_PATH}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>${LOG_DIR}/agent.log</string>
    <key>StandardErrorPath</key>
    <string>${LOG_DIR}/agent-error.log</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
PLIST

chmod 644 "$PLIST_PATH"

# Remove quarantine flag so unsigned builds launch without Gatekeeper blocking.
xattr -rd com.apple.quarantine "$APP_PATH" 2>/dev/null || true

# Copy baked-in config from the .app bundle to system config dir if present.
BAKED_CONFIG_DIR="${APP_PATH}/Contents/Resources/config"
BAKED_ENV="${BAKED_CONFIG_DIR}/agent.env"
BAKED_JSON="${BAKED_CONFIG_DIR}/collector.json"
if [ -f "$BAKED_ENV" ]; then
    cp "$BAKED_ENV" "$APP_SUPPORT_DIR/agent.env"
    chmod 644 "$APP_SUPPORT_DIR/agent.env"
    echo "  Pre-configured agent.env installed to $APP_SUPPORT_DIR/"
fi
if [ -f "$BAKED_JSON" ]; then
    cp "$BAKED_JSON" "$APP_SUPPORT_DIR/collector.json"
    chmod 644 "$APP_SUPPORT_DIR/collector.json"
    echo "  Pre-configured collector.json installed to $APP_SUPPORT_DIR/"
fi

# Load the LaunchDaemon (system domain)
launchctl load "$PLIST_PATH" 2>/dev/null || true

echo "Detec Agent: post-install complete."
echo "  LaunchDaemon installed: $PLIST_PATH (starts at boot)."
echo "  Config: $APP_SUPPORT_DIR"
if [ -f "$BAKED_ENV" ] || [ -f "$APP_SUPPORT_DIR/agent.env" ]; then
    echo "  Server config found. The agent will connect automatically."
else
    echo "  Add agent.env or collector.json to $APP_SUPPORT_DIR and restart: sudo launchctl kickstart -k system/com.detec.agent"
fi
exit 0
